# root/etc/services.d/librevault/run
#!/usr/bin/with-contenv bash

exec \
	librevault-daemon
